# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/banip.list
/www/cgi-bin/banip
/etc/config/banip
/etc/banip/banip.allowlist
/usr/lib/banip-functions.sh
/etc/init.d/banip
/usr/bin/banip-service.sh
/etc/banip/banip.blocklist
/etc/banip/banip.countries
/etc/banip/banip.tpl
/etc/banip/banip.feeds
/etc/banip/banip.custom.feeds
