# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/pbr.list
/etc/hotplug.d/iface/70-pbr
/usr/share/nftables.d/chain-post/mangle_prerouting/30-pbr.nft
/etc/config/pbr
/usr/share/nftables.d/chain-post/mangle_postrouting/30-pbr.nft
/usr/share/nftables.d/chain-post/mangle_forward/30-pbr.nft
/usr/share/pbr/pbr.firewall.include
/etc/init.d/pbr
/etc/uci-defaults/90-pbr
/usr/share/pbr/pbr.user.aws
/usr/share/nftables.d/chain-post/mangle_input/30-pbr.nft
/usr/share/nftables.d/table-post/30-pbr.nft
/usr/share/nftables.d/chain-post/mangle_output/30-pbr.nft
/usr/share/pbr/pbr.user.netflix
