# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/smcroute.list
/usr/bin/smcroutectl
/usr/sbin/smcrouted
/etc/smcroute.conf
/etc/init.d/smcroute
