# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/stubby.list
/etc/stubby/stubby.yml.default
/etc/stubby/stubby.yml
/usr/sbin/stubby
/etc/init.d/stubby
/etc/config/stubby
