# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/tailscaled.list
/etc/config/tailscale
/usr/sbin/tailscaled
/lib/upgrade/keep.d/tailscaled
/etc/init.d/tailscale
