# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/usbmuxd.list
/usr/sbin/usbmuxd
/lib/upgrade/keep.d/usbmuxd
/etc/init.d/usbmuxd
