# buildEnv/openwrt-imagebuilder-23.05.2-rockchip-armv8.Linux-x86_64/build_dir/target-aarch64_generic_musl/root-rockchip/usr/lib/opkg/info/watchcat.list
/etc/config/watchcat
/etc/uci-defaults/migrate-watchcat
/usr/bin/watchcat.sh
/etc/init.d/watchcat
