# buildEnv/openwrt-imagebuilder-23.05.2-ramips-mt7621.Linux-x86_64/build_dir/target-mipsel_24kc_musl/root.orig-ramips/usr/lib/opkg/info/wpad-mesh-mbedtls.list
/usr/sbin/wpad
/usr/sbin/hostapd
/etc/init.d/radius
/usr/sbin/hostapd-radius
/usr/share/hostap/wpa_supplicant.uc
/etc/radius/clients
/lib/upgrade/keep.d/wpad-mesh-mbedtls
/etc/config/radius
/usr/share/hostap/hostapd.uc
/etc/radius/users
/usr/sbin/wpa_supplicant
